# tazpup-core-files/pre-slitaz/var/lib/tazpkg/installed/hal/files.list
/etc/dbus-1/system.d/hal.conf
/etc/init.d/hald
/usr/share/hal/fdi/policy/10osvendor/10-input-policy.fdi
/usr/share/hal/fdi/policy/10osvendor/10-power-mgmt-policy.fdi
/usr/share/hal/fdi/policy/10osvendor/10-laptop-panel-mgmt-policy.fdi
/usr/share/hal/fdi/policy/10osvendor/10-rfkill-switch.fdi
/usr/share/hal/fdi/policy/10osvendor/10-tabletPCs.fdi
/usr/share/hal/fdi/policy/10osvendor/10-x11-input.fdi
/usr/share/hal/fdi/policy/10osvendor/10-leds.fdi
/usr/share/hal/fdi/policy/10osvendor/30-ntfs-write.fdi
/usr/share/hal/fdi/policy/10osvendor/20-storage-methods.fdi
/usr/share/hal/fdi/fdi.dtd
/usr/bin/hal-disable-polling
/usr/bin/lshal
/usr/bin/hal-set-property
/usr/bin/hal-device
/usr/bin/hal-get-property
/usr/bin/hal-find-by-capability
/usr/bin/hal-find-by-property
/usr/bin/hal-is-caller-locked-out
/usr/bin/hal-lock
/usr/bin/hald
/usr/lib/libhal.so.1.0.0
/usr/lib/hal/hal-system-power-pmu
/usr/lib/hal/hald-addon-input
/usr/lib/hal/hal-storage-cleanup-all-mountpoints
/usr/lib/hal/hald-probe-storage
/usr/lib/hal/hald-addon-leds
/usr/lib/hal/hal-system-power-pm-is-supported
/usr/lib/hal/hal-system-sonypic
/usr/lib/hal/hald-addon-rfkill-killswitch
/usr/lib/hal/hald-addon-generic-backlight
/usr/lib/hal/hal-storage-mount
/usr/lib/hal/hald-probe-input
/usr/lib/hal/scripts/hal-system-wol-enabled
/usr/lib/hal/scripts/hal-luks-setup
/usr/lib/hal/scripts/hal-system-killswitch-get-power
/usr/lib/hal/scripts/hal-system-wol-supported
/usr/lib/hal/scripts/hal-dockstation-undock
/usr/lib/hal/scripts/hal-luks-teardown
/usr/lib/hal/scripts/hal-system-power-suspend
/usr/lib/hal/scripts/hal-system-power-reboot
/usr/lib/hal/scripts/hal-system-power-shutdown
/usr/lib/hal/scripts/hal-system-killswitch-set-power
/usr/lib/hal/scripts/hal-functions
/usr/lib/hal/scripts/hal-system-power-suspend-hybrid
/usr/lib/hal/scripts/hal-system-lcd-set-brightness
/usr/lib/hal/scripts/hal-system-lcd-get-brightness
/usr/lib/hal/scripts/linux/hal-system-power-suspend-hybrid-linux
/usr/lib/hal/scripts/linux/hal-system-power-suspend-linux
/usr/lib/hal/scripts/linux/hal-system-power-set-power-save-linux
/usr/lib/hal/scripts/linux/hal-luks-remove-linux
/usr/lib/hal/scripts/linux/hal-system-wol-linux
/usr/lib/hal/scripts/linux/hal-system-lcd-set-brightness-linux
/usr/lib/hal/scripts/linux/hal-system-wol-enabled-linux
/usr/lib/hal/scripts/linux/hal-luks-teardown-linux
/usr/lib/hal/scripts/linux/hal-system-wol-enable-linux
/usr/lib/hal/scripts/linux/hal-system-power-reboot-linux
/usr/lib/hal/scripts/linux/hal-system-killswitch-set-power-linux
/usr/lib/hal/scripts/linux/hal-system-lcd-get-brightness-linux
/usr/lib/hal/scripts/linux/hal-system-power-hibernate-linux
/usr/lib/hal/scripts/linux/hal-system-power-shutdown-linux
/usr/lib/hal/scripts/linux/hal-luks-setup-linux
/usr/lib/hal/scripts/linux/hal-system-killswitch-get-power-linux
/usr/lib/hal/scripts/linux/hal-system-wol-supported-linux
/usr/lib/hal/scripts/linux/hal-dockstation-undock-linux
/usr/lib/hal/scripts/hal-system-power-hibernate
/usr/lib/hal/scripts/hal-luks-remove
/usr/lib/hal/scripts/hal-system-power-set-power-save
/usr/lib/hal/scripts/hal-system-wol-enable
/usr/lib/hal/hald-runner
/usr/lib/hal/hal-storage-closetray
/usr/lib/hal/hal-storage-eject
/usr/lib/hal/hald-addon-acpi
/usr/lib/hal/hald-generate-fdi-cache
/usr/lib/hal/hal-system-setserial
/usr/lib/hal/hal-storage-unmount
/usr/lib/hal/hald-probe-smbios
/usr/lib/hal/hald-probe-serial
/usr/lib/hal/hald-addon-pmu
/usr/lib/hal/hald-probe-volume
/usr/lib/hal/hald-probe-hiddev
/usr/lib/hal/hald-addon-storage
/usr/lib/hal/hal-storage-cleanup-mountpoint
/usr/lib/libhal-storage.so.1.0.0
/usr/lib/libhal.so.1
/usr/lib/libhal.so
/usr/lib/libhal-storage.so.1
/usr/lib/libhal-storage.so
