# tazpup-core-files/pre-slitaz/var/lib/tazpkg/installed/slitaz-base-files/files.list
/lib/libtaz.sh
/var/www/style.css
/var/www/cgi-bin/cgi-env.sh
/var/www/cgi-bin/index.cgi
/var/www/index.html
/var/log/wtmp
/sbin/mktazdevs.sh
/etc/nsswitch.conf
/etc/daemons.conf
/etc/profile
/etc/sysctl.conf
/etc/group
/etc/init.d/rc.functions
/etc/fstab
/etc/gshadow
/etc/slitaz-release
/etc/adduser.conf
/etc/slitaz/slitaz.conf
/etc/hosts
/etc/shadow
/etc/networks
/etc/host.conf
/etc/passwd
/etc/motd
/etc/hostname
/etc/shells
/etc/services
/etc/issue
/etc/securetty
/usr/lib/slitaz/libpkg.sh
/usr/lib/slitaz/httphelper.sh
/usr/share/doc/slitaz-doc.js
/usr/share/doc/slitaz/httphelper.txt
/usr/share/doc/slitaz/libpkg.txt
/usr/share/doc/slitaz/libtaz.txt
/usr/share/doc/slitaz-doc.css
/usr/share/locale/ja/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/de/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/vi/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/ca/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/pt_BR/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/zh_TW/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/it/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/fa/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/el/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/fr/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/pl/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/zh_CN/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/ru/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/sv/LC_MESSAGES/slitaz-base.mo
/usr/share/locale/es/LC_MESSAGES/slitaz-base.mo
/usr/share/applications/httpd-cgi.desktop
/usr/bin/ldd
/usr/bin/ipinfo
/usr/bin/calc
/usr/bin/slitaz
/usr/bin/sudo
/usr/bin/man
/dev/stdin
/dev/vcsa
/dev/vcs
/dev/fd
/dev/core
/dev/stdout
/dev/flash
/dev/ram
/dev/stderr
/etc/mtab
/usr/lib/slitaz/httphelper
