# tazpup-core-files/slitaz/etc/xdg/autostart/check-acpi.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
NoDisplay=true
Name=check-acpi
Exec=check-acpi
Name[en_US.UTF-8]=check-acpi
