# tazpup-core-files/slitaz/etc/xdg/autostart/cursor-fix.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
NoDisplay=true
Name=cursor-fix
Exec=cursorfixer
Name[en_US.UTF-8]=cursor-fix
