# tazpup-core-files/slitaz/etc/xdg/autostart/pupevent.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
NoDisplay=true
Name=Pup Event Frontend
Exec=pup_event_frontend_d
