# tazpup-core-files/pre-slitaz/var/lib/tazpkg/installed/busybox/files.list
/var/www/httpd/404.html
/etc/busybox.conf
/etc/zcip.script
/etc/dnsd.conf
/etc/init.d/daemon
/etc/inetd.conf
/etc/httpd.conf
/etc/syslog.conf
/etc/resolv.conf
/etc/udhcpd.conf
/bin/busybox
/usr/share/udhcpc/default.script
/usr/share/applications/httpd.desktop
/usr/share/applications/vi.desktop
/sbin/iproute
/sbin/vconfig
/sbin/route
/sbin/watchdog
/sbin/iplink
/sbin/blockdev
/sbin/findfs
/sbin/syslogd
/sbin/poweroff
/sbin/shutdown
/sbin/fbsplash
/sbin/arp
/sbin/adjtimex
/sbin/getty
/sbin/iprule
/sbin/ipneigh
/sbin/raidautorun
/sbin/mkdosfs
/sbin/init
/sbin/reboot
/sbin/loadkmap
/sbin/swapoff
/sbin/swapon
/sbin/tc
/sbin/nameif
/sbin/mkfs.vfat
/sbin/mkswap
/sbin/ip
/sbin/switch_root
/sbin/sysctl
/sbin/setconsole
/sbin/freeramdisk
/sbin/klogd
/sbin/fdisk
/sbin/blkid
/sbin/iptunnel
/sbin/sulogin
/sbin/pivot_root
/sbin/losetup
/sbin/tunctl
/sbin/hdparm
/sbin/halt
/sbin/fstrim
/sbin/zcip
/sbin/ipaddr
/sbin/ifconfig
/sbin/hwclock
/sbin/udhcpc
/sbin/ifenslave
/etc/init.d/crond
/etc/init.d/ntpd
/etc/init.d/syslogd
/etc/init.d/httpd
/etc/init.d/lpd
/etc/init.d/tftpd
/etc/init.d/udhcpd
/etc/init.d/inetd
/etc/init.d/ftpd
/etc/init.d/klogd
/etc/init.d/dnsd
/etc/init.d/zcip
/bin/fgrep
/bin/dmesg
/bin/sh
/bin/df
/bin/su
/bin/false
/bin/cp
/bin/vi
/bin/tar
/bin/ps
/bin/rpm
/bin/fatattr
/bin/gunzip
/bin/sleep
/bin/fsync
/bin/mktemp
/bin/grep
/bin/login
/bin/netstat
/bin/pidof
/bin/true
/bin/fdflush
/bin/base64
/bin/cat
/bin/setserial
/bin/echo
/bin/printenv
/bin/pwd
/bin/mknod
/bin/kbd_mode
/bin/uname
/bin/nice
/bin/stat
/bin/sed
/bin/ln
/bin/date
/bin/chmod
/bin/touch
/bin/mt
/bin/zcat
/bin/sync
/bin/mpstat
/bin/ip
/bin/setpriv
/bin/scriptreplay
/bin/ping
/bin/pipe_progress
/bin/setarch
/bin/dd
/bin/usleep
/bin/iostat
/bin/kill
/bin/ls
/bin/ping6
/bin/dumpkmap
/bin/conspy
/bin/arch
/bin/egrep
/bin/gzip
/bin/ipcalc
/bin/chgrp
/bin/fbvnc
/bin/cpio
/bin/hostname
/bin/mkdir
/bin/chown
/bin/rm
/bin/rev
/bin/ionice
/bin/ash
/bin/more
/bin/rmdir
/bin/stty
/bin/mv
/bin/dnsdomainname
/bin/watch
/bin/uncompress
/usr/sbin/rdev
/usr/sbin/chroot
/usr/sbin/brctl
/usr/sbin/ether-wake
/usr/sbin/crond
/usr/sbin/ntpd
/usr/sbin/fakeidentd
/usr/sbin/httpd
/usr/sbin/lpd
/usr/sbin/deluser
/usr/sbin/rtcwake
/usr/sbin/arping
/usr/sbin/chpasswd
/usr/sbin/tftpd
/usr/sbin/udhcpd
/usr/sbin/rdate
/usr/sbin/inetd
/usr/sbin/partprobe
/usr/sbin/addgroup
/usr/sbin/fdformat
/usr/sbin/i2cdetect
/usr/sbin/ftpd
/usr/sbin/killall5
/usr/sbin/dhcprelay
/usr/sbin/nologin
/usr/sbin/ubidetach
/usr/sbin/fsfreeze
/usr/sbin/setlogcons
/usr/sbin/nbd-client
/usr/sbin/powertop
/usr/sbin/dnsd
/usr/sbin/delgroup
/usr/sbin/adduser
/usr/sbin/readprofile
/usr/sbin/sendmail
/usr/sbin/ubiattach
/usr/sbin/fbset
/usr/bin/printf
/usr/bin/sum
/usr/bin/lspci
/usr/bin/tr
/usr/bin/chrt
/usr/bin/microcom
/usr/bin/time
/usr/bin/clear
/usr/bin/du
/usr/bin/sha3sum
/usr/bin/lsof
/usr/bin/ftpget
/usr/bin/ftpput
/usr/bin/id
/usr/bin/w
/usr/bin/whois
/usr/bin/volname
/usr/bin/unix2dos
/usr/bin/yes
/usr/bin/lpr
/usr/bin/mkfifo
/usr/bin/env
/usr/bin/head
/usr/bin/eject
/usr/bin/traceroute
/usr/bin/crontab
/usr/bin/ipcrm
/usr/bin/unzip
/usr/bin/split
/usr/bin/httpd_helper.sh
/usr/bin/hostid
/usr/bin/expr
/usr/bin/chvt
/usr/bin/diff
/usr/bin/logname
/usr/bin/tftp
/usr/bin/ttysize
/usr/bin/sort
/usr/bin/unexpand
/usr/bin/bunzip2
/usr/bin/cal
/usr/bin/lzcat
/usr/bin/basename
/usr/bin/fallocate
/usr/bin/flock
/usr/bin/logger
/usr/bin/mkpasswd
/usr/bin/reset
/usr/bin/realpath
/usr/bin/nohup
/usr/bin/tac
/usr/bin/setfattr
/usr/bin/vlock
/usr/bin/shred
/usr/bin/strings
/usr/bin/md5sum
/usr/bin/users
/usr/bin/uuencode
/usr/bin/find
/usr/bin/traceroute6
/usr/bin/setsid
/usr/bin/tty
/usr/bin/dirname
/usr/bin/rpm2cpio
/usr/bin/which
/usr/bin/who
/usr/bin/wget
/usr/bin/dpkg-deb
/usr/bin/whoami
/usr/bin/telnet
/usr/bin/openvt
/usr/bin/hd
/usr/bin/last
/usr/bin/lsusb
/usr/bin/pmap
/usr/bin/script
/usr/bin/nmeter
/usr/bin/wall
/usr/bin/hexdump
/usr/bin/od
/usr/bin/ipcs
/usr/bin/renice
/usr/bin/mesg
/usr/bin/seq
/usr/bin/lpq
/usr/bin/groups
/usr/bin/fold
/usr/bin/deallocvt
/usr/bin/killall
/usr/bin/unxz
/usr/bin/fuser
/usr/bin/beep
/usr/bin/passwd
/usr/bin/unlzma
/usr/bin/timeout
/usr/bin/free
/usr/bin/expand
/usr/bin/patch
/usr/bin/sha1sum
/usr/bin/bzcat
/usr/bin/test
/usr/bin/uptime
/usr/bin/[
/usr/bin/comm
/usr/bin/tee
/usr/bin/xzcat
/usr/bin/pstree
/usr/bin/cmp
/usr/bin/readlink
/usr/bin/uniq
/usr/bin/nslookup
/usr/bin/sha256sum
/usr/bin/taskset
/usr/bin/wc
/usr/bin/pkill
/usr/bin/tail
/usr/bin/pscan
/usr/bin/pgrep
/usr/bin/nc
/usr/bin/uudecode
/usr/bin/install
/usr/bin/cut
/usr/bin/xargs
/usr/bin/dumpleases
/usr/bin/udhcpc6
/usr/bin/top
/usr/bin/cksum
/usr/bin/sha512sum
/usr/bin/[[
/usr/bin/cryptpw
/usr/bin/dos2unix
/usr/bin/awk
/usr/bin/smemcap
/usr/bin/less
/usr/bin/setkeycodes
