# tazpup-core-files/pre-slitaz/var/lib/tazpkg/installed/slim/files.list
/etc/slim.conf
/etc/init.d/slim
/usr/share/slim/strings
/usr/share/slim/themes/base/panel.png
/usr/share/slim/themes/base/background.png
/usr/share/slim/themes/base/slim.theme
/usr/bin/slim
/usr/bin/slim-theme
