# tazpup-core-files/slitaz/etc/udev/rules.d/99-usb-tethering.rules
ACTION=="add", SUBSYSTEM=="net", ENV{ID_BUS}=="usb", ENV{ID_USB_DRIVER}=="rndis_host", RUN+="/usr/bin/net-interface $env{INTERFACE} start"
