# tazpup-core-files/pre-slitaz/var/lib/tazpkg/installed/tazpanel/files.list
/var/www/tazpanel/lib/libtazpanel
/var/www/tazpanel/lib/tazpanel.js.gz
/var/www/tazpanel/README
/var/www/tazpanel/settings.cgi
/var/www/tazpanel/doc/tazpanel.ru.html
/var/www/tazpanel/doc/tazpanel.de.html
/var/www/tazpanel/doc/tazpanel.en.html
/var/www/tazpanel/doc/tazpanel.fr.html
/var/www/tazpanel/doc/tazpanel.pt.html
/var/www/tazpanel/hardware.cgi
/var/www/tazpanel/boot.cgi
/var/www/tazpanel/styles/default/footer.html
/var/www/tazpanel/styles/default/header.html
/var/www/tazpanel/styles/default/tazpanel.css.gz
/var/www/tazpanel/styles/default/tazpanel.ru.html
/var/www/tazpanel/styles/default/tazpanel.de.html
/var/www/tazpanel/styles/default/tazpanel.en.html
/var/www/tazpanel/styles/default/favicon.ico
/var/www/tazpanel/styles/default/tazpanel.fr.html
/var/www/tazpanel/styles/default/tazpanel.pt.html
/var/www/tazpanel/styles/default/tazpanel.ttf
/var/www/tazpanel/network.cgi
/var/www/tazpanel/README.html
/var/www/tazpanel/hosts.cgi
/var/www/tazpanel/index.cgi
/var/www/tazpanel/help.cgi
/var/log/tazpanel.log
/etc/slitaz/tazpanel.conf
/etc/slitaz/httpd.conf
/usr/share/icons/hicolor/32x32/apps/tazpanel.png
/usr/share/icons/hicolor/48x48/apps/tazpanel.png
/usr/share/icons/hicolor/16x16/apps/tazpanel.png
/usr/share/locale/ja/LC_MESSAGES/tazpanel.mo
/usr/share/locale/de/LC_MESSAGES/tazpanel.mo
/usr/share/locale/pt_BR/LC_MESSAGES/tazpanel.mo
/usr/share/locale/vi_VN/LC_MESSAGES/tazpanel.mo
/usr/share/locale/it/LC_MESSAGES/tazpanel.mo
/usr/share/locale/el/LC_MESSAGES/tazpanel.mo
/usr/share/locale/fr/LC_MESSAGES/tazpanel.mo
/usr/share/locale/pl/LC_MESSAGES/tazpanel.mo
/usr/share/locale/zh_CN/LC_MESSAGES/tazpanel.mo
/usr/share/locale/ru/LC_MESSAGES/tazpanel.mo
/usr/share/locale/sv/LC_MESSAGES/tazpanel.mo
/usr/share/locale/es/LC_MESSAGES/tazpanel.mo
/usr/share/applications/tazpanel.desktop
/usr/bin/tazpanel
/var/www/tazpanel/user
/var/www/tazpanel/doc/tazpanel.html
/etc/init.d/tazpanel
/usr/share/fonts/TTF/tazpanel.ttf
