# tazpup-core-files/pre-slitaz/var/lib/tazpkg/installed/slitaz-tools/files.list
/sbin/tazkeymap
/sbin/hwsetup
/sbin/tazlocale
/sbin/slitaz-config
/sbin/tazhw
/etc/locale.conf
/etc/init.d/firewall
/etc/slitaz/firewall.conf
/etc/slitaz/firewall.sh
/etc/firewall.conf
/etc/TZ
/etc/keymap.conf
/usr/share/doc/slitaz/post-install.txt
/usr/share/licenses/COPYING
/usr/share/licenses/mit.txt
/usr/share/licenses/lgpl.txt
/usr/share/licenses/mozilla.txt
/usr/share/licenses/gpl.txt
/usr/share/licenses/bsd.txt
/usr/share/locale/ja/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/de/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/pt_BR/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/it/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/fa/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/el/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/fr/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/pl/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/zh_CN/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/ru/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/sv/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/es/LC_MESSAGES/slitaz-tools.mo
/usr/sbin/setmixer
/usr/bin/editor
/usr/bin/startd
/usr/bin/file-manager
/usr/bin/history
/usr/bin/tazx
/usr/bin/browser
/usr/bin/stopd
/usr/bin/decode
/usr/bin/terminal
/usr/bin/startx
/usr/bin/sdt
/usr/bin/frugal
