# tazpup-core-files/pre-slitaz/var/lib/tazpkg/installed/acpid/files.list
/etc/init.d/acpid
/etc/acpi/acpi_handler.sh
/etc/acpi/events/ac_adapter
/etc/acpi/events/buttons
/usr/bin/acpi_listen
/usr/sbin/acpid
/usr/sbin/kacpimon
